# lighttpd.upstart
description "Lighttpd Web Server"

start on pld.network-started
stop on pld.shutdown-started

respawn

console output

env SHELL=/bin/sh

script
	# Source old-style service configuration
	[ -f /etc/sysconfig/lighttpd ] && . /etc/sysconfig/lighttpd
	exec /usr/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf $HTTPD_OPTS
end script
